# backend/app/services/agent_bundle/opspilot-agent.service
[Unit]
Description=OpsPilot metrics push agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/opspilot
Environment=OPSPILOT_AGENT_CONFIG=/opt/opspilot/config/agent.json
ExecStart=/usr/bin/python3 /opt/opspilot/bin/opspilot-agent.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
